# auth_service/Supervisord/supervisord.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php artisan schedule:work
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/log/supervisor/laravel-scheduler.err.log
stdout_logfile=/var/log/supervisor/laravel-scheduler.out.log
user=root

[program:command-app-kafka-consume]
process_name=%(program_name)s_%(process_num)02d
command=php artisan app:kafka:consume
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/log/supervisor/command-app-kafka-consume.err.log
stdout_logfile=/var/log/supervisor/command-app-kafka-consume.out.log
user=root
